# systemd_user/.config/systemd/user/packages.service
[Unit]
Description=Update Package list

[Service]
Type=oneshot
ExecStart=/bin/bash /home/akreuzer/.config/scripts/update_package_list.sh
