# systemd_user/.config/systemd/user/packages.timer
[Unit]
Description=test

[Timer]
OnUnitActiveSec=1d
OnBootSec=15min

[Install]
WantedBy=timers.target
