# systemd_root/root/zfs/sanoid/sanoid.timer
[Unit]
Description=sanoid zfs snapshoting timer

[Timer]
OnCalendar=minutely
OnBootSec=5min

[Install]
WantedBy=timers.target
