# systemd_root/root/backups/backblaze.timer
[Unit]
Description=backblaze sync

[Timer]
OnCalendar=daily
OnBootSec=15min

[Install]
WantedBy=timers.target
