# systemd/root/backups/backblaze.service
[Service] 
Type=simple 
ExecStart=/usr/local/bin/backblaze.sh
