# systemd/root/zfs/sanoid/sanoid.service
[Unit]
Description=sanoid zfs snapshoting

[Service] 
Type=simple 
ExecStart=/usr/bin/sanoid --cron
