# systemd_root/root/zfs/daily/zfs-snapshots-daily.service
[Unit]
Description=zfs daily snapshot
[Service] 
Type=simple 
ExecStart=/usr/local/bin/zfs-snapshots-daily.sh
